# disk/etc/inittab

rc1:1:wait:/bin/klogd
rcl:1:wait:/etc/rc.local

0:1:respawn:/bin/getty -t /dev/tty0
S1:1:respawn:/bin/getty -t /dev/ttyS1


